# deploy/tsuyu-rag.service
[Unit]
Description=TSUYU RAG Chatbot
After=network.target postgresql.service ollama.service
Wants=postgresql.service ollama.service

[Service]
Type=simple
User=tsuyu
WorkingDirectory=/opt/tsuyu-rag
# Muat .env dari laluan mutlak (penting — service tak warisi shell env)
EnvironmentFile=/opt/tsuyu-rag/.env
ExecStart=/opt/tsuyu-rag/tsuyu-rag-chatbot
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
